# systemd/voxi.service
[Unit]
Description=Voxi Voice Dictation Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/voxi/.venv/bin/python -m src.daemon
Restart=on-failure
RestartSec=3
WorkingDirectory=%h/.local/share/voxi/

[Install]
WantedBy=default.target
